# Distribution/PackageRoot/Library/LaunchDaemons/com.gmaendel.SiemensIPGiver.bpf.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.gmaendel.SiemensIPGiver.bpf</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Application Support/SiemensIPGiver/SiemensIPGiverBPFSetup.sh</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/var/log/SiemensIPGiverBPFSetup.err</string>
	<key>StandardOutPath</key>
	<string>/var/log/SiemensIPGiverBPFSetup.log</string>
</dict>
</plist>
